# deploy/systemd/nxms-mailbox.service
[Unit]
Description=NXMS Mailbox Service
After=network-online.target tor.service
Wants=network-online.target

[Service]
Type=simple
User=nxms
Group=nxms
UMask=0077
WorkingDirectory=/opt/freeforum/nexum/escrow
EnvironmentFile=-/etc/nxms-mailbox.env
ExecStart=/opt/freeforum/nexum/escrow/nxms-mailbox/target/release/nxms-mailbox serve
Restart=always
RestartSec=3
MemoryMax=384M
TasksMax=256
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
ReadWritePaths=/var/lib/nxms

[Install]
WantedBy=multi-user.target
